# crontab
* * * * * . $HOME/.backey.env; /usr/bin/backey
